# deploy/chainer-trainer.service
[Unit]
Description=Chainer Agent - PPO Training Service (PyTorch GPU)
After=network.target

[Service]
Type=simple
User=m
WorkingDirectory=/home/m/chainer-agent
Environment=TRAINER_PORT=5555
ExecStart=/home/m/chainer-agent/venv/bin/python3 training/trainer.py
Restart=always
RestartSec=5
StandardOutput=append:/home/m/chainer-agent/trainer.log
StandardError=append:/home/m/chainer-agent/trainer.log

[Install]
WantedBy=multi-user.target
